# systemd/slack-app.service
[Unit]
Description=Slack App

[Service]
ExecStart=/opt/slack-app/scripts/start.sh

[Install]
WantedBy=multi-user.target
